# ops/config/systemd/nginx-iris.socket
[Unit]
Description=Socket for iris nginx proxy

[Socket]
ListenStream=0.0.0.0:16649

[Install]
WantedBy=sockets.target
